# docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
username=supervisor
password=supervisor

[supervisord]
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=supervisor
password=supervisor

[program:api]
; Run migrations then start uvicorn. Priority 10 = starts before nginx.
command=sh -c "cd /app && alembic upgrade head && exec uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4"
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
